# etc/init.d/tlp-intel-mobile-openrc
#!/sbin/openrc-run

description="Modify Nitrux TLP settings for mobile Intel processors based on processor generation"

depend() {
    need tlp-openrc
    after bootmisc
}

cpu_model_check() {
    local model
    model=$(grep -m1 "model name" /proc/cpuinfo 2>/dev/null | sed -n 's/.*i[0-9]-\([0-9]\{3,5\}[A-Z]*\).*/\1/p')

    if [ -z "$model" ]; then
        echo "not_intel"
    elif echo "$model" | grep -Eq '^[0-9]{3,5}[A-Z]?$'; then
        echo "$model"
    else
        echo "unknown"
    fi
}

cpu_generation_check() {
    local model generation
    model=$(cpu_model_check)

    case "$model" in
        not_intel)
            echo "not_intel"
            return 0
            ;;
        unknown)
            echo "unknown"
            return 1
            ;;
        *H)
            generation=${model%%[A-Z]*}
            if [ ${#generation} -eq 4 ]; then
                generation=$((generation / 1000))
            elif [ ${#generation} -eq 5 ]; then
                generation=$((generation / 1000))
            fi
            echo "$generation"
            ;;
        *)
            echo "not_mobile"
            ;;
    esac
}

modify_tlp_settings() {
    local tlp_file="/etc/tlp.d/99-nx-tlp-bat-settings.conf"
    local generation
    generation=$(cpu_generation_check)

    if [ "$generation" = "not_intel" ]; then
        ewarn "Non-Intel processor detected. No changes made to $tlp_file."
        return 0
    fi

    if [ "$generation" = "not_mobile" ]; then
        ewarn "Non-mobile Intel processor detected. No changes made to $tlp_file."
        return 0
    fi

    if [ "$generation" = "unknown" ]; then
        eerror "Unable to determine CPU generation. No changes made to $tlp_file."
        return 1
    fi

    if [ ! -f "$tlp_file" ]; then
        eerror "$tlp_file not found. Unable to apply changes."
        return 1
    fi

    if [ "$generation" -ge 11 ]; then
        sed -i 's/^CPU_SCALING_GOVERNOR_ON_BAT=.*/CPU_SCALING_GOVERNOR_ON_BAT=userspace/' "$tlp_file"
        einfo "Set scaling governor to 'userspace' for generation $generation mobile processors."
    else
        sed -i 's/^CPU_SCALING_GOVERNOR_ON_BAT=.*/CPU_SCALING_GOVERNOR_ON_BAT=ondemand/' "$tlp_file"
        einfo "Set scaling governor to 'ondemand' for generation $generation mobile processors."
    fi
}

start() {
    ebegin "Adjusting Nitrux TLP settings for Intel mobile CPUs based on processor generation"
    modify_tlp_settings
    eend $?
}
