# etc/init.d/tlp-intel-mobile-openrc
#!/sbin/openrc-run

description="Modify Nitrux TLP settings for mobile Intel processors based on processor generation"

depend() {
    need tlp-openrc
    after bootmisc
}

cpu_model_check() {
    if grep -q "Intel" /proc/cpuinfo; then
        model=$(grep -m1 "model name" /proc/cpuinfo | sed -n 's/.* \([0-9]\{4\}[A-Z]\{1,2\}\)$/\1/p')
        if echo "$model" | grep -Eq "^[0-9]{4}[A-Z]{1,2}$"; then
            generation=$(echo "$model" | cut -c1-2)
            suffix=$(echo "$model" | cut -c5-)
            case "$suffix" in
                H|HQ|HK|HF|Y|U|B)
                    echo "$generation"
                    ;;
                *)
                    echo "not_mobile"
                    ;;
            esac
        else
            echo "unknown"
        fi
    else
        echo "not_intel"
    fi
}

modify_tlp_settings() {
    local tlp_file="/etc/tlp.d/99-nx-bat-settings.conf"
    local generation
    generation=$(cpu_model_check)

    if [ "$generation" = "not_intel" ]; then
        ewarn "Non-Intel processor detected. No changes made to $tlp_file."
        return 0
    fi

    if [ "$generation" = "not_mobile" ]; then
        ewarn "Non-mobile Intel processor detected. No changes made to $tlp_file."
        return 0
    fi

    if [ "$generation" = "unknown" ]; then
        ewarn "Unable to determine CPU generation. No changes made to $tlp_file."
        return 1
    fi

    if [ "$generation" -ge 11 ]; then
        sed -i 's/^CPU_SCALING_GOVERNOR_ON_BAT=.*/CPU_SCALING_GOVERNOR_ON_BAT=userspace/' "$tlp_file"
        einfo "Set scaling governor to 'userspace' for generation $generation mobile processors."
    else
        sed -i 's/^CPU_SCALING_GOVERNOR_ON_BAT=.*/CPU_SCALING_GOVERNOR_ON_BAT=ondemand/' "$tlp_file"
        einfo "Set scaling governor to 'ondemand' for generation $generation mobile processors."
    fi
}

start() {
    ebegin "Adjusting Nitrux TLP settings for Intel mobile CPUs based on processor generation"
    modify_tlp_settings
    eend $?
}
