# etc/init.d/tlp-intel-mobile-openrc
#!/sbin/openrc-run

description="Modify Nitrux TLP settings for mobile Intel processors based on processor generation"

depend() {
    need tlp-openrc
    after bootmisc
}

cpu_model_check() {
    model=$(grep -m1 "model name" /proc/cpuinfo 2>/dev/null | sed -n 's/.*i[0-9]-\([0-9]\{3,5\}[A-Z]*\).*/\1/p')

    if echo "$model" | grep -Eq '^[0-9]{3,5}[A-Z]?$'; then
        echo "$model"
    else
        echo "unknown"
    fi
}

cpu_generation_check() {
    local model=$(cpu_model_check)
    
    case "$model" in
        unknown)
            echo "unknown"
            ;;
        [0-9][0-9][0-9][A-Z]?)
            echo "${model:0:1}" # Extract generation from 3-digit model (e.g., 9750H -> 9)
            ;;
        [0-9][0-9][0-9][0-9][A-Z]?)
            echo "${model:0:2}" # Extract generation from 4-digit model (e.g., 13620H -> 13)
            ;;
        *)
            echo "not_intel"
            ;;
    esac
}

modify_tlp_settings() {
    local tlp_file="/etc/tlp.d/99-nx-bat-settings.conf"
    local generation=$(cpu_generation_check)

    case "$generation" in
        unknown)
            ewarn "Unable to determine CPU generation. No changes made to $tlp_file."
            return 1
            ;;
        not_intel)
            ewarn "Non-Intel processor detected. No changes made to $tlp_file."
            return 0
            ;;
        1[1-9]|[2-9][0-9]) # 11th gen or newer
            sed -i 's/^CPU_SCALING_GOVERNOR_ON_BAT=.*/CPU_SCALING_GOVERNOR_ON_BAT=userspace/' "$tlp_file"
            einfo "Set scaling governor to 'userspace' for generation $generation mobile processors."
            ;;
        [6-9]) # 6th to 10th gen
            sed -i 's/^CPU_SCALING_GOVERNOR_ON_BAT=.*/CPU_SCALING_GOVERNOR_ON_BAT=ondemand/' "$tlp_file"
            einfo "Set scaling governor to 'ondemand' for generation $generation mobile processors."
            ;;
        *)
            ewarn "Unrecognized CPU generation: $generation. No changes made to $tlp_file."
            return 1
            ;;
    esac
}

start() {
    ebegin "Adjusting Nitrux TLP settings for Intel mobile CPUs based on processor generation"
    modify_tlp_settings
    eend $?
}
