# debian/postinst
#!/usr/bin/env bash

# -- Add service to runlevel.

rc-update add tlp-openrc default
rc-update add tlp-intel-mobile-openrc default

# -- Update services list.

rc-update -u
